# nightscout-docker/rootfs/etc/services.d/nightscout/run
#!/usr/bin/execlineb -P
node "/app/server.js"